# supervisors/cron_schedule.ini
[program:xxxx-计划任务]
command=php think cron:schedule
directory=$站点目录$
autorestart=true
startsecs=3
startretries=3
stdout_logfile=$站点目录$/cron_schedule.out.log
stderr_logfile=$站点目录$/cron_schedule.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=www
priority=999
numprocs=1
process_name=%(program_name)s_%(process_num)02d
